# probe/node/systemd/neox2-realm-probe-node.service
[Unit]
Description=neoX2 Realm Probe (Node.js)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-/etc/neox2-realm-probe.env
WorkingDirectory=/opt/neox2-realm-probe-node
ExecStart=/usr/bin/node /opt/neox2-realm-probe-node/server.js
Restart=on-failure
RestartSec=3
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target

